# openclaw/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:dashboard]
command=python -c "from polymarket_mcp.web.app import start; start()"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:openclaw]
command=script -qec "node /usr/lib/node_modules/openclaw/dist/index.js gateway --allow-unconfigured --port 18789 --bind lan --verbose 2>&1" /dev/null
directory=/root/.openclaw
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=10
startsecs=5
priority=2
environment=HOME="/root",OPENCLAW_STATE_DIR="/root/.openclaw",OPENCLAW_NO_RESPAWN="1",TERM="xterm-256color"
